# files/etc/supervisord.d/dovecot.ini
[program:dovecot]
command=/usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F
autostart=true
autorestart=true

stdout_syslog=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
# Redirect output so we can see it using "docker logs"
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
