# files/etc/supervisord.d/postfix.ini
[program:postfix]
process_name=master
directory=/etc/postfix
command=/usr/sbin/postfix start
startsecs=0
autorestart=false

stdout_syslog=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
# Redirect output so we can see it using "docker logs"
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
